# bonescript/suite/jessie/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

echo "NODE_PATH=/usr/local/lib/node_modules" > /etc/default/node
echo "export NODE_PATH=/usr/local/lib/node_modules" > /etc/profile.d/node.sh
chmod 755 /etc/profile.d/node.sh

systemctl enable bonescript.socket || true
systemctl enable bonescript-autorun.service || true

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

