# bb-wl18xx-firmware/suite/bullseye/debian/bb-wl18xx-softap.bb-wl18xx-softap.service
[Unit]
Description=BeagleBoard.org wl18xx SoftAp0 Service
After=usb-gadget.target
ConditionFileIsExecutable=/usr/bin/bb-wl18xx-softap

[Service]
Type=oneshot
ExecStartPre=/sbin/modprobe wlcore
ExecStartPre=/sbin/modprobe wlcore_sdio
ExecStartPre=/sbin/modprobe wl18xx
ExecStart=/usr/bin/bb-wl18xx-softap
Type=forking

[Install]
WantedBy=multi-user.target
