# bb-wl18xx-firmware/suite/bullseye/debian/prerm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

# avoid cancelling of "stop" via socket or D-Bus activation
if [ -d /run/systemd/system ] && [ "$1" = remove ]; then
	systemctl mask --runtime bb-wl18xx-bluetooth.service || true
	systemctl mask --runtime bb-wl18xx-wlan0.service || true
	# stop ${service}.service manually since invoke-rc.d currently doesn't
	# deal with masked services (see Bug #717109)
	systemctl stop bb-wl18xx-bluetooth.service || true
	systemctl stop bb-wl18xx-wlan0.service || true

	if [ -f /lib/systemd/system/bb-wl18xx-tether.service ] ; then
		systemctl mask --runtime bb-wl18xx-tether.service || true
		systemctl stop bb-wl18xx-tether.service || true
	fi
fi

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    remove|upgrade|deconfigure)
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
