# bela-customizations/suite/bullseye/debian/ssh_shutdown.service
[Unit]
Description=Kill ssh sessions before network
After=network.target

[Service]
TimeoutStartSec=0
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/bin/killall sshd

[Install]
WantedBy=multi-user.target
