# bb-code-server/suite/jammy/debian/user/code-server.service
[Unit]
Description=code-server
Wants=network.target dphys-swapfile.service
After=network.target

[Service]
Type=exec
ExecStart=/opt/bb-code-server/app/bin/code-server --disable-telemetry --disable-update-check
Restart=always

[Install]
WantedBy=default.target
