# bb-wl18xx-firmware/suite/stretch/debian/bb-wl18xx-wlan0.service
[Unit]
Description=BB WL18xx wlan0 Service
After=basic.target

[Service]
ExecStart=/usr/bin/bb-wl18xx-wlan0

[Install]
WantedBy=multi-user.target

