# debian-8-bb-node-red-installer/suite/jessie/debian/node-red.socket
[Socket]
ListenStream=1880

[Install]
WantedBy=sockets.target
