# bb-bonescript-installer-beta/suite/stretch/debian/bonescript-autorun.service
[Unit]
Description=Bonescript autorun
ConditionPathExists=|/var/lib/cloud9

[Service]
WorkingDirectory=/usr/local/lib/node_modules/bonescript
EnvironmentFile=/etc/default/node
ExecStart=/usr/bin/nodejs autorun.js
SyslogIdentifier=bonescript-autorun

[Install]
WantedBy=multi-user.target
