# bb-wlan0-defaults/suite/focal/debian/81-add-mesh0-interface.rules
SUBSYSTEM=="ieee80211", ACTION=="add", KERNEL=="phy[0-9]", DEVPATH=="*mmc/*", \
	RUN+="/sbin/iw phy %k interface add mesh0 type mp"
