# bb-code-server/suite/buster/debian/bb-code-server.socket
[Unit]
Description=bb-code-server Socket
PartOf=bb-code-server.service

[Socket]
ListenStream=3000

[Install]
WantedBy=sockets.target
