# bb-code-server/suite/buster/debian/bb-code-server.bb-code-server.service
[Unit]
Description=BeagleBoard.org Code Server
Wants=network.target
ConditionFileIsExecutable=/opt/bb-code-server/app/bin/code-server
After=multi-user.target

[Service]
Type=exec
ExecStart=/opt/bb-code-server/app/bin/code-server --auth none --disable-telemetry --disable-update-check --config /opt/bb-code-server/config/config.yaml --extensions-dir /opt/bb-code-server/extensions/ --bind-addr 0.0.0.0:3000
Restart=always
User=debian

[Install]
WantedBy=multi-user.target
