# bb-code-server/suite/buster/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

package="bb-code-server"
service="bb-code-server"

main_install () {
	mkdir -p /var/lib/bb-code-server/config/
	mkdir -p /var/lib/bb-code-server/extensions/

	/opt/bb-code-server/app/bin/code-server --config /var/lib/bb-code-server/config/config.yaml --extensions-dir /var/lib/code-server/extensions/ --list-extensions
	/opt/bb-code-server/app/bin/code-server --config /var/lib/bb-code-server/config/config.yaml --extensions-dir /var/lib/code-server/extensions/ --install-extension formulahendry.code-runner
	/opt/bb-code-server/app/bin/code-server --config /var/lib/bb-code-server/config/config.yaml --extensions-dir /var/lib/code-server/extensions/ --install-extension donjayamanne.githistory
	/opt/bb-code-server/app/bin/code-server --config /var/lib/bb-code-server/config/config.yaml --extensions-dir /var/lib/code-server/extensions/ --list-extensions

	chown -R 1000:1000 /var/lib/bb-code-server/config/
	chown -R 1000:1000 /var/lib/bb-code-server/extensions/

	echo "#!/bin/bash" > /opt/code-server/start.sh
	echo "/opt/bb-code-server/app/bin/code-server --auth none --disable-telemetry --disable-update-check --config /var/lib/bb-code-server/config/config.yaml --extensions-dir /var/lib/bb-code-server/extensions/ --bind-addr 0.0.0.0:3000" >> /opt/bb-code-server/start.sh
	chmod +x /opt/bb-code-server/start.sh

	systemctl disable ${service}.socket || true
	systemctl disable ${service} || true
	#systemctl enable ${service}.socket || true
	#systemctl start ${service} || true
	echo "${package}:Installed"
}

#cleanup () {
#}

case "$1" in
    configure)
	main_install
#	cleanup
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

