# c9-core-installer/suite/buster/debian/cloud9.socket
[Unit]
Description=Cloud9 Socket
PartOf=cloud9.service

[Socket]
ListenStream=3000

[Install]
WantedBy=sockets.target
