# rcnee-internal-testlab/suite/bullseye/debian/postinst
#!/bin/sh

set -e

if [ -f /etc/systemd/timesyncd.conf ] ; then
	sed -i -e 's:#NTP=:NTP=192.168.1.24:g' /etc/systemd/timesyncd.conf
fi

if [ -f /lib/systemd/system/apt-daily.service ] ; then
	systemctl disable apt-daily.service || true
	systemctl disable apt-daily.timer || true
fi

if [ -f /lib/systemd/system/apt-daily-upgrade.service ] ; then
	systemctl disable apt-daily-upgrade.service || true
	systemctl disable apt-daily-upgrade.timer || true
fi

if [ -f /etc/apt/apt.conf.d/02-no-pdiffs ] ; then
	rm -f /etc/apt/apt.conf.d/02-no-pdiffs || true
fi

if [ -f /etc/apt/apt.conf.d/02-no-languages ] ; then
	rm -f /etc/apt/apt.conf.d/02-no-languages || true
fi

if [ -f /lib/systemd/system/bb-code-server.service ] ; then
	systemctl enable bb-code-server.service || true
fi


new_machine=$(cat /proc/device-tree/model | sed "s/ /_/g" | sed "s/-/_/g" | tr -d '\000')
case "${new_machine}" in
BeagleBoard.org_BeagleBone_AI)
	if [ -f /etc/udev/rules.d/81-add-SoftAp0-interface.rules ] ; then
		rm -rf /etc/udev/rules.d/81-add-SoftAp0-interface.rules || true
	fi
	if [ -f /etc/udev/rules.d/82-SoftAp0-start-hostpad.rules ] ; then
		rm -rf /etc/udev/rules.d/82-SoftAp0-start-hostpad.rules || true
	fi
	;;
esac

#DEBHELPER#
