# bb-code-server/suite/bullseye/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

package="bb-code-server"
service="bb-code-server"
wdir="/var/lib/bb-code-server"
wfile="/opt/bb-code-server/start.sh"

main_install () {
	mkdir -p ${wdir}/config/
	mkdir -p ${wdir}/extensions/

	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --list-extensions

	#Git
	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --install-extension donjayamanne.githistory --force

	#Code Runner
	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --install-extension formulahendry.code-runner --force

	#Python
	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --install-extension ms-toolsai.jupyter --force
	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --install-extension ms-python.python --force

	#bug: https://github.com/cdr/code-server/issues/1448
#platformio.platformio-ide doesn't seem to actually work in the vscode browser..
#	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --install-extension ms-vscode.cpptools --force
#	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --install-extension platformio.platformio-ide --force

	/opt/bb-code-server/app/bin/code-server --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --list-extensions

	chown -R 1000:1000 ${wdir}/config/
	chown -R 1000:1000 ${wdir}/extensions/

	echo "#!/bin/bash" > ${wfile}
	echo "/opt/bb-code-server/app/bin/code-server --auth none --disable-telemetry --disable-update-check --config ${wdir}/config/config.yaml --extensions-dir ${wdir}/extensions/ --bind-addr 0.0.0.0:3000" >> ${wfile}
	chmod +x ${wfile}

	systemctl disable ${service}.socket || true
	systemctl disable ${service} || true
	#systemctl enable ${service}.socket || true
	#systemctl start ${service} || true
	echo "${package}:Installed"
}

#cleanup () {
#}

case "$1" in
    configure)
	main_install
#	cleanup
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

