# bela-customizations/suite/bullseye/debian/bela-customizations/lib/systemd/system/bela_startup.service
[Unit]
Description=Run Bela at boot
DefaultDependencies=no
Requires=systemd-modules-load.service
Conflicts=shutdown.target
Before=shutdown.target

[Service]
User=root
EnvironmentFile=/opt/Bela/startup_env
ExecStart=/usr/bin/stdbuf -oL -eL /opt/Bela/bela_startup.sh
Type=simple
Restart=always

[Install]
WantedBy=default.target
