# debian-9-bluealsa/suite/stretch/debian/bluealsa.postinst
#!/bin/sh

set -e

case "${1}" in
  configure)
   if dpkg --compare-versions "${2}" lt "0.8" && deb-systemd-helper --quiet is-enabled bluealsa.service; then
     deb-systemd-helper disable bluealsa.service
   fi
   ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

  *)
    echo "postinst called with unknown argument \`${1}'" >&2
    exit 1
    ;;
esac

#DEBHELPER#
