# bb-customizations/suite/trixie/debian/prerm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

remove_file () {
	if [ -f ${wfile} ] ; then
		rm -f ${wfile} || true
	fi
}

remove_dir () {
	if [ -d ${wdir} ] ; then
		rm -rf ${wdir} || true
	fi
}

remove_all () {
	wfile="/etc/udev/rules.d/uio.rules" ; remove_file
	if [ -f /lib/systemd/system/generic-board-startup.service ] ; then
		systemctl stop generic-board-startup.service || true
		systemctl disable generic-board-startup.service || true
		wfile="/lib/systemd/system/generic-board-startup.service" ; remove_file
	fi
}

case "$1" in
    remove|upgrade|deconfigure)
	remove_all
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
