# generic-sys-mods/suite/bullseye/debian/generic-sys-mods.bbbio-set-sysconf.service
[Unit]
Description=BeagleBoard Set up system configuration
Requires=boot-firmware.mount
After=boot-firmware.mount

[Service]
Type=oneshot
ExecStart=/usr/sbin/bbbio-set-sysconf

[Install]
RequiredBy=basic.target
