# mpm/suite/stretch/debian/mpm.upstart
description "Multi Proc Manager Daemon"
author "Hao Zhang <hzhang@ti.com>"
version "0.1"
start on runlevel [2345] 
stop on runlevel [!2345] 
exec /usr/sbin/mpmsrv

