# xenomai/suite/bullseye/debian/libxenomai1/DEBIAN/prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit/13.3.4
if [ -x "/etc/init.d/xenomai" ] && [ "$1" = remove ]; then
	invoke-rc.d xenomai stop || exit 1
fi
# End automatically added section
