# dk-hackaday-2021/suite/bullseye/debian/99-hd44780.rules
SUBSYSTEM=="hd44780", ACTION=="add", \
	RUN+="/bin/chgrp gpio /dev/lcd0", \
	RUN+="/bin/chmod g=u /dev/lcd0"

