# debian-8-connman/suite/xenial/debian/connman.upstart
description "Connection Manager"

start on started dbus
stop on stopping dbus

expect fork
respawn

exec connmand --nobacktrace
