# bbb.io-beagleconnect/suite/buster/debian/beagleconnect.lowpan.service
[Unit]
Description=Create 6lowpan (IEEE802.15.4) network device
After=network.target

[Service]
EnvironmentFile=/etc/default/lowpan
Type=oneshot
User=root
ExecStart=/usr/sbin/create_lowpan $CHN $PAN $IP6
ExecStop=/usr/sbin/delete_lowpan
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
