# bb-wlan0-defaults/suite/bullseye/debian/81-add-SoftAp0-interface.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="wlan[0-9]", DRIVERS=="wl18xx_driver" \
	RUN+="/sbin/iw dev %k interface add SoftAp0 type managed"
