# bb-wlan0-defaults/suite/bullseye/debian/81-ap-wl18xx-interface.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="wlan0", DRIVERS=="wl18xx_driver" \
	RUN+="/sbin/iw dev %k interface add ap%n type managed"
