# bb-doc-bone101-jekyll/suite/jessie/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

package="bb-doc-bone101-jekyll"
service="jekyll-autorun"

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

git_clone () {
	mkdir -p ${git_target_dir} || true
	git clone ${git_repo} ${git_target_dir} --depth 1 || true
}

pkg_install () {
	git_repo="https://github.com/beagleboard/bone101"
	git_target_dir="/var/lib/cloud9"
	git_clone

	if [ -f ${git_target_dir}/.git/config ] ; then
		cd ${git_target_dir}/

		echo "jekyll pre-building bone101"
		/usr/bin/jekyll build --destination bone101
	fi

	systemctl enable ${service}.service || true
	systemctl start ${service} || true
	echo "${package}:Installed"
}

case "$1" in
    configure)
	pkg_install
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

