# bb-wl18xx-firmware/suite/focal/debian/bb-wl18xx-wlan0.service
[Unit]
Description=BB WL18xx wlan0 Service
After=local-fs.target

[Service]
ExecStart=/usr/bin/bb-wl18xx-wlan0
Type=forking
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

