# bb-wl18xx-firmware/suite/focal/debian/bb-wl18xx-bluetooth.service
[Unit]
Description=BB WL18xx Bluetooth Service
After=getty.target bluetooth.service generic-board-startup.service 

[Service]
ExecStart=/usr/bin/bb-wl18xx-bluetooth
Type=forking

[Install]
WantedBy=multi-user.target

