# seeed-wificonfig-installer/suite/jessie/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

package="seeed-wificonfig-installer"
service="wifidog-gateway"

npm_pkg_install () {
	echo "${package}:installing:[${npm_file}]"
	if [ ! -d /usr/local/lib/node_modules/ ] ; then
		mkdir -p /usr/local/lib/node_modules/ || true
	fi
	if [ -d /usr/local/lib/node_modules/${npm_project}/ ] ; then
		rm -rf /usr/local/lib/node_modules/${npm_project}/ || true
	fi

	echo "${package}:updating:[/usr/local/lib/node_modules/${npm_project}/]"
	tar xf /var/cache/${package}/${npm_file} -C /usr/local/lib/node_modules/
}

npm_install () {
	node_bin="/usr/bin/nodejs"

	unset node_version
	node_version=$(/usr/bin/nodejs --version || true)
	echo "${package}:nodejs:[${node_version}]"

	v12="v0.12.17"
	v4="v4.6.2"
	v6="v6.9.1"

	wificonfig="wificonfig-0.1.5-c59f06b"

	npm_project="wificonfig"
	case "${node_version}" in
	v0.12.*)
		npm_file="${wificonfig}-${v12}.tar.xz"
		npm_pkg_install
		;;
	v4.*)
		npm_file="${wificonfig}-${v4}.tar.xz"
		npm_pkg_install
		;;
	v6.*)
		npm_file="${wificonfig}-${v6}.tar.xz"
		npm_pkg_install
		;;
	esac

	if [ ! -d ${NODE_PATH}/${pkg}/log/ ] ; then
		mkdir -p ${NODE_PATH}/${pkg}/log/ || true
	fi

	systemctl disable wificonfig.socket || true
	systemctl disable wificonfig.service || true

	#disable bonescript port 80...
	systemctl stop bonescript.socket || true
	systemctl disable bonescript.socket || true

	systemctl stop bonescript-autorun.service || true
	systemctl disable bonescript-autorun.service || true
	mkdir -p /etc/wificonfig/ || true
	echo "${package}:Installed"
}

case "$1" in
    configure)
	systemctl stop wificonfig.service || true
	npm_install
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

