# seeed-wificonfig-installer/suite/jessie/debian/prerm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

service="wifidog-gateway"

# avoid cancelling of "stop" via socket or D-Bus activation
if [ -d /run/systemd/system ] && [ "$1" = remove ]; then
    systemctl mask --runtime ${service}.service
    systemctl mask --runtime wificonfig.socket
    # stop ${service}.service manually since invoke-rc.d currently doesn't
    # deal with masked services (see Bug #717109)
    systemctl stop ${service}.service || true
    systemctl stop wificonfig.socket || true
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
