# seeed-wificonfig-installer/suite/jessie/debian/prerm
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

# avoid cancelling of "stop" via socket or D-Bus activation
if [ -d /run/systemd/system ] && [ "$1" = remove ]; then
	systemctl mask --runtime wifidog-pre-startup.service || true
	systemctl mask --runtime wifidog-gateway.service || true
	systemctl mask --runtime wificonfig.socket || true
	# stop ${service}.service manually since invoke-rc.d currently doesn't
	# deal with masked services (see Bug #717109)
	systemctl stop wifidog-gateway.service || true
	systemctl stop wificonfig.service || true

	#removed:
	if [ -f /lib/systemd/system/wifidog-pre-startup.service ] ; then
		systemctl stop wifidog-pre-startup.service || true
		systemctl disable wifidog-pre-startup.service || true
	fi
fi

package="seeed-wificonfig-installer"

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

remove_file () {
	if [ -f /var/cache/${package}/${wfile} ] ; then
		rm -f /var/cache/${package}/${wfile} || true
	fi
}

remove_dir () {
	if [ -d ${wdir} ] ; then
		rm -rf ${wdir} || true
	fi
}

remove_all() {
	wfile="wificonfig-0.1.0-e4e30e7-v0.12.13.tar.xz" ; remove_file
	wfile="wificonfig-0.1.4-bc0ac5e-v0.12.13.tar.xz" ; remove_file
	wfile="wificonfig-0.1.5-8da2e95-v0.12.13.tar.xz" ; remove_file

	#Always keep the current one for quick reinstalls...
	#wfile="wificonfig-0.1.5-2a7ba99-v0.12.15-bbgw-bone101-dir-change.tar.xz" ; remove_file
}

case "$1" in
    remove|upgrade|deconfigure)
	#remove_all
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
