# seeed-wificonfig-installer/suite/jessie/debian/wifidog-gateway.service
[Unit]
Description=wifidog gateway Service
After=wifidog-pre-startup.service

[Service]
Type=simple
ExecStart=/usr/bin/wifidog -s
KillSignal=SIGINT
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
