# seeed-wificonfig-installer/suite/jessie/debian/wificonfig.service
[Unit]
Description=wificonfig Service

[Service]
WorkingDirectory=/usr/local/lib/node_modules/wificonfig
ExecStart=/usr/bin/nodejs /usr/bin/npm start
SyslogIdentifier=wificonfig
