# seeed-wificonfig-installer/suite/jessie/debian/wificonfig.socket
[Socket]
ListenStream=520

[Install]
WantedBy=sockets.target
