# bb-wl18xx-firmware/suite/yakkety/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

echo_default () {
	echo "# TETHER_ENABLED: Whether or not to run the /usr/bin/bb-wl18xx-tether daemon; set to no to disable." > /etc/default/bb-wl18xx
	echo "TETHER_ENABLED=yes" >> /etc/default/bb-wl18xx
	echo "" >> /etc/default/bb-wl18xx
	echo "# USE_WL18XX_IP_PREFIX: default IP block of SoftAP0 interface" >> /etc/default/bb-wl18xx
	echo "USE_WL18XX_IP_PREFIX=\"192.168.8\"" >> /etc/default/bb-wl18xx
	echo "" >> /etc/default/bb-wl18xx
	echo "# USE_GENERATED_DNSMASQ: use generated version of /etc/dnsmasq.d/SoftAp0; set to no so user can modify /etc/dnsmasq.d/SoftAp0" >> /etc/default/bb-wl18xx
	echo "USE_GENERATED_DNSMASQ=yes" >> /etc/default/bb-wl18xx
	echo "" >> /etc/default/bb-wl18xx
	echo "# USE_GENERATED_HOSTAPD: use generated version of /etc/hostapd.conf; set to no so user can modify /etc/hostapd.conf" >> /etc/default/bb-wl18xx
	echo "USE_GENERATED_HOSTAPD=yes" >> /etc/default/bb-wl18xx
	echo "" >> /etc/default/bb-wl18xx
	echo "# USE_PERSONAL_SSID: set custom ssid (overides BeagleBone/BeagleBone-WXYZ)" >> /etc/default/bb-wl18xx
	echo "#USE_PERSONAL_SSID=" >> /etc/default/bb-wl18xx
	echo "" >> /etc/default/bb-wl18xx
}

if [ ! -f /etc/default/bb-wl18xx ] ; then
	echo_default
else
	unset test_new_option
	test_new_option=$(cat /etc/default/bb-wl18xx | grep TETHER_ENABLED || true)
	if [ "x${test_new_option}" = "x" ] ; then
		echo_default
	fi

	unset test_new_option
	test_new_option=$(cat /etc/default/bb-wl18xx | grep USE_GENERATED_DNSMASQ || true)
	if [ "x${test_new_option}" = "x" ] ; then
		echo "" >> /etc/default/bb-wl18xx
		echo "# USE_GENERATED_DNSMASQ: use generated version of /etc/dnsmasq.d/SoftAp0; set to no so user can modify /etc/dnsmasq.d/SoftAp0" >> /etc/default/bb-wl18xx
		echo "USE_GENERATED_DNSMASQ=yes" >> /etc/default/bb-wl18xx
		echo "" >> /etc/default/bb-wl18xx
		echo "# USE_GENERATED_HOSTAPD: use generated version of /etc/hostapd.conf; set to no so user can modify /etc/hostapd.conf" >> /etc/default/bb-wl18xx
		echo "USE_GENERATED_HOSTAPD=yes" >> /etc/default/bb-wl18xx
		echo "" >> /etc/default/bb-wl18xx
	fi

	unset test_new_option
	test_new_option=$(cat /etc/default/bb-wl18xx | grep USE_WL18XX_IP_PREFIX || true)
	if [ "x${test_new_option}" = "x" ] ; then
		echo "" >> /etc/default/bb-wl18xx
		echo "# USE_WL18XX_IP_PREFIX: default IP block of SoftAP0 interface" >> /etc/default/bb-wl18xx
		echo "USE_WL18XX_IP_PREFIX=\"192.168.8\"" >> /etc/default/bb-wl18xx
		echo "" >> /etc/default/bb-wl18xx
	fi

	unset test_new_option
	test_new_option=$(cat /etc/default/bb-wl18xx | grep USE_PERSONAL_SSID || true)
	if [ "x${test_new_option}" = "x" ] ; then
		echo "" >> /etc/default/bb-wl18xx
		echo "# USE_PERSONAL_SSID: set custom ssid (overides BeagleBone/BeagleBone-WXYZ)" >> /etc/default/bb-wl18xx
		echo "#USE_PERSONAL_SSID=" >> /etc/default/bb-wl18xx
		echo "" >> /etc/default/bb-wl18xx
	fi
fi

case "$1" in
	configure)
		systemctl enable bb-wl18xx-bluetooth.service || true
		systemctl enable bb-wl18xx-wlan0.service || true
	;;

	abort-upgrade|abort-remove|abort-deconfigure)
	;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

