# debian-9-c9-core-installer/suite/stretch/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

package="c9-core-installer"

core_version="3.1.4056.git20171215"
#c9-core_3.1.4056.git20171215-v6.12.2-build.tar.xz

#deb https://deb.nodesource.com/node_6.x jessie main
   wfile_6_x="c9-core_${core_version}-v6.12.2-build.tar.xz"

extract_dir="/opt/cloud9/build/"

export NODE_PATH=/usr/local/lib/node_modules

npm_install () {
	if [ -d /etc/avahi/ ] ; then
		#Annouce http server via DNS Sevice Discovery
		wfile="/etc/avahi/services/cloud9.service"
		echo "<?xml version=\"1.0\" standalone='no'?><!--*-nxml-*-->" > ${wfile}
		echo "<!DOCTYPE service-group SYSTEM \"avahi-service.dtd\">" >> ${wfile}
		echo "" >> ${wfile}
		echo "<!-- See avahi.service(5) for more information about this configuration file -->" >> ${wfile}
		echo "" >> ${wfile}
		echo "<service-group>" >> ${wfile}
		echo "" >> ${wfile}
		echo "  <name replace-wildcards=\"yes\">Cloud9 IDE for %h</name>" >> ${wfile}
		echo "  <service>" >> ${wfile}
		echo "" >> ${wfile}
		echo "    <type>_http._tcp</type>" >> ${wfile}
		echo "    <port>3000</port>" >> ${wfile}
		echo "  </service>" >> ${wfile}
		echo "" >> ${wfile}
		echo "</service-group>" >> ${wfile}
	fi
}

dl_install() {
	unset node_version
	node_version=$(/usr/bin/nodejs --version || true)
	case "${node_version}" in
	v6.*)
		wfile=${wfile_6_x}
		;;
	*)
		wfile=${wfile_6_x}
		;;
	esac

	cd /var/cache/${package}/

	systemctl daemon-reload || true
	echo "${package}:Stopping:cloud9.service"
	systemctl stop cloud9.service || true
	echo "${package}:Stopping:cloud9.socket"
	systemctl stop cloud9.socket || true

	if [ ! -d ${extract_dir} ] ; then
		mkdir -p ${extract_dir} || true
	fi

	echo "${package}:extracting:${wfile}"
	tar xf ${wfile} -C ${extract_dir}
	chown -R 1000:1000 ${extract_dir}

	if [ ! -d /var/lib/cloud9 ] ; then
		mkdir -p /var/lib/cloud9 || true
	fi

	systemctl daemon-reload || true
}

c9_directory () {
	if [ ! -d /opt/cloud9/.c9/node/bin/ ] ; then
		mkdir -p /opt/cloud9/.c9/node/bin/ || true
	fi
	ln -sf /usr/bin/node /opt/cloud9/.c9/node/bin/node

	groupadd --system cloud9ide || true
	chown -R root:cloud9ide /opt/cloud9/.c9/ || true
	chmod -R g+w /opt/cloud9/.c9/ || true

	mkdir -p /opt/cloud9/.cache/ || true
	chown -R root:cloud9ide  /opt/cloud9/.cache/ || true
	chmod -R g+w  /opt/cloud9/.cache/ || true

	systemctl enable cloud9.socket || true
	systemctl restart cloud9.socket || true
	echo "${package}:Installed"
}

case "$1" in
    configure)
	npm_install
	dl_install
	c9_directory
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

