# tmux-1.6/debian/jessie/debian/preinst
#!/bin/sh

set -e

# Versions before 1.5-2 did not properly disable the tmux-cleanup init script.
if [ "$1" = "upgrade" ] && dpkg --compare-versions "$2" 'lt' "1.5-2"; then
    update-rc.d -f tmux-cleanup remove >/dev/null
fi

#DEBHELPER#
