# generic-sys-mods/suite/bookworm/debian/61-smsc95xx-mac-address.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="eth0", DRIVERS=="smsc95xx" \
	RUN+="/usr/sbin/bb-smsc95xx-set-mac-address"
