# generic-sys-mods/suite/bookworm/debian/generic-sys-mods.grow_partition.service
[Unit]
Description=Grow Root Partition
Before=ssh.service
ConditionFileIsExecutable=/usr/bin/bb-growpart

[Service]
Type=oneshot
ExecStart=/usr/bin/bb-growpart
ExecStartPost=systemctl enable resize_filesystem
ExecStartPost=systemctl disable grow_partition

[Install]
WantedBy=multi-user.target
