# gcp-iot-demo/suite/jessie/debian/gcp-iot-demo.socket
[Socket]
ListenStream=3001

[Install]
WantedBy=sockets.target
