# gcp-iot-demo/suite/jessie/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

set -e

package="gcp-iot-demo"
service="gcp-iot-demo"

npm_pkg_install () {
	echo "${package}:installing:[${npm_file}]"
	if [ ! -d /usr/local/lib/node_modules/ ] ; then
		mkdir -p /usr/local/lib/node_modules/ || true
	fi
	if [ -d /usr/local/lib/node_modules/${npm_project}/ ] ; then
		rm -rf /usr/local/lib/node_modules/${npm_project}/ || true
	fi

	echo "${package}:updating:[/usr/local/lib/node_modules/${npm_project}/]"
	tar xf /var/cache/${package}/${npm_file} -C /usr/local/lib/node_modules/
}

npm_install () {
	node_bin="/usr/bin/nodejs"

	unset node_version
	node_version=$(/usr/bin/nodejs --version || true)
	echo "${package}:nodejs:[${node_version}]"

	v12="v0.12.17"
	v4="v4.6.2"
	v6="v6.9.1"

	node_gcp="gcp-iot-demo-0.0.1-fb406df"

	npm_project="gcp-iot-demo"
	case "${node_version}" in
	v0.12.*)
		npm_file="${node_gcp}-${v12}.tar.xz"
		npm_pkg_install
		;;
	v4.*)
		npm_file="${node_gcp}-${v4}.tar.xz"
		npm_pkg_install
		;;
	v6.*)
		npm_file="${node_gcp}-${v6}.tar.xz"
		npm_pkg_install
		;;
	esac

	if [ -d /etc/avahi/services/ ] ; then
		#Annouce http server via DNS Sevice Discovery
		wfile="/etc/avahi/services/gcp-iot-demo.service"
		echo "<?xml version=\"1.0\" standalone='no'?><!--*-nxml-*-->" > ${wfile}
		echo "<!DOCTYPE service-group SYSTEM \"avahi-service.dtd\">" >> ${wfile}
		echo "" >> ${wfile}
		echo "<!-- See avahi.service(5) for more information about this configuration file -->" >> ${wfile}
		echo "" >> ${wfile}
		echo "<service-group>" >> ${wfile}
		echo "" >> ${wfile}
		echo "  <name replace-wildcards=\"yes\">IoT developer kit icebreaker for %h</name>" >> ${wfile}
		echo "  <service>" >> ${wfile}
		echo "" >> ${wfile}
		echo "    <type>_http._tcp</type>" >> ${wfile}
		echo "    <port>3001</port>" >> ${wfile}
		echo "  </service>" >> ${wfile}
		echo "" >> ${wfile}
		echo "</service-group>" >> ${wfile}
		chown -R root:root ${wfile}
	fi

	systemctl enable ${service}.socket || true
	systemctl start ${service}.socket || true
	echo "${package}:Installed"
}

npm_install

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

