# bb-wl18xx-firmware/suite/zesty/debian/bb-wl18xx-tether.service
[Unit]
Description=BB WL18xx tether Service

[Service]
ExecStart=/usr/bin/bb-wl18xx-tether
Type=forking

[Install]
WantedBy=multi-user.target

