# bb-usb-gadgets/suite/bullseye/debian/bb-usb-gadgets.run_bb_usb_gadgets.service
[Unit]
Description=Run USB gadgets
Before=console-setup.service
ConditionFileIsExecutable=/usr/bin/run-bb-usb-gadgets

[Service]
Type=oneshot
ExecStart=/usr/bin/run-bb-usb-gadgets

[Install]
WantedBy=multi-user.target
