# bb-node-red-installer/suite/jessie/debian/node-red.service
[Unit]
Description=Start Node-RED

[Service]
Environment=HOME=/root
WorkingDirectory=/root
ExecStart=/usr/bin/node-red.sh
SyslogIdentifier=Node-RED
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
